# templates/default/supervisord.exhibitor.conf.erb
[program:exhibitor]
command = java -jar /var/exhibitor/target/exhibitor-1.5.5.jar -c file
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/exhibitor.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/exhibitor.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:exhibitor_server]
programs=exhibitor
