# templates/default/supervisord.zookeeper.conf.erb
[program:zookeeper]
command = /var/zookeeper-<%=@version%>/bin/zkServer.sh start-foreground
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/zookeeper.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/zookeeper.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:zookeeper_server]
programs=zookeeper
