# templates/default/supervisord.zookeeper.conf.erb
[program:zookeeper]
command = <%=Chef::Config[:file_cache_path]%>/zookeeper-<%=node['zookeeper']['version']%>/bin/zkServer.sh start-foreground
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/zookeeper.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/zookeeper.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[group:zookeeper_server]
programs=zookeeper
